# mk-52.service
[Unit]
Description=МК-52 emulator controller
# Wait for filesystem; we don't use the network.
After=local-fs.target

[Service]
# Reference values — tools/install-pi.sh auto-detects and rewrites these.
# Edit ExecStart, WorkingDirectory, and User to match your install if you
# install by hand. Use pypy3 instead of python3 for ~65× chip-loop speedup.
ExecStart=/usr/bin/python3 -u app.py
WorkingDirectory=/home/pi/mk-52-retrofit/controller

Restart=always
RestartSec=5

# Send stdout/stderr to journald. View live with:
#   journalctl -u mk-52 -f
StandardOutput=journal
StandardError=journal

User=pi

[Install]
WantedBy=multi-user.target
